# host/deploy/peartune-host.service
[Unit]
Description=PearTune host (self-hosted music library over P2P)
Documentation=https://github.com/peerloomllc/peartune
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# The PearTune host is plain Node - it runs host/index.js directly (unlike the
# PearCircle seeder, there is no Bare worklet to supervise). Node 20+ must be on
# PATH. Stage the repo at /opt/peartune and install the host's deps ONCE with:
#   cd /opt/peartune/host && npm ci --omit=dev
# (the host resolves ../protocol and ../client, so keep the whole repo, not just
# host/). Config comes from the EnvironmentFile - see peartune-host.env.example.
User=peartune
Group=peartune
WorkingDirectory=/opt/peartune
EnvironmentFile=/etc/peartune/peartune-host.env
ExecStart=/usr/bin/env node /opt/peartune/host/index.js
Restart=always
RestartSec=5

# Hardening. The host only needs to WRITE its data dir; everything else is
# read-only or off. ReadWritePaths must match PEARTUNE_DATA in the env file.
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
ReadWritePaths=/var/lib/peartune

[Install]
WantedBy=multi-user.target
